# infra/ansible/roles/observability_stack/templates/restic-restore-test.timer.j2
# {{ ansible_managed }}
[Unit]
Description=BagTrip weekly Restic restore drill

[Timer]
OnCalendar={{ observability_restic_restore_test_schedule }}
Persistent=true
RandomizedDelaySec=30min
Unit=restic-restore-test.service

[Install]
WantedBy=timers.target
